# deploy/link-forge.service
[Unit]
Description=Link Forge - Personal Knowledge Forager
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=musicofhel
WorkingDirectory=/home/musicofhel/link-forge
ExecStart=/home/musicofhel/.nvm/versions/node/v22.22.0/bin/node dist/index.js
Restart=on-failure
RestartSec=10
TimeoutStopSec=30
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
